# usr/share/mate/autostart/caja-dropbox-autostart.desktop
[Desktop Entry]
Name=Caja Dropbox autostart
Comment=If Caja Dropbox is installed, make sure it autostarts.
Exec=/usr/lib/ubuntu-mate/caja-dropbox-autostart
Terminal=false
Type=Application
Icon=caja-dropbox
Categories=Network;FileTransfer;
StartupNotify=false
OnlyShowIn=MATE;
NoDisplay=true
